# for build/ecat-bot.conf
[program:ecat-bot]
directory=/home/ecatalog
command=python3 tgbot.py
autostart=true
autorestart=true
stderr_logfile=/var/log/ecatalog.err.log
stdout_logfile=/var/log/ecatalog.out.log