# for build/ecat-monitor.conf
[program:ecat-monitor]
directory=/home/ecatalog
command=python3 monitor.py
autostart=true
autorestart=true
stderr_logfile=/var/log/ecatalogm.err.log
stdout_logfile=/var/log/ecatalogm.out.log