# for build/mp-monitor.conf
[program:mp-monitor]
directory=/home/multyparse
command=python3 monitor.py
autostart=true
autorestart=true
stderr_logfile=/var/log/mpmonitor.err.log
stdout_logfile=/var/log/mpmonitor.out.log